# user-src/systemd/wg-fips-quick.target
[Unit]
Description=WolfGuard Tunnels via wg-fips-quick(8)
